# operations/supervisord.d/mqlservd.ini
[program:mcbridgefsd]
process_name = %(program_name)s_%(process_num)02d
command = /usr/local/bin/mcbridgefsd
autostart = true
autorestart = true
startretries = 50
user = gtarcea
numprocs = 1
redirect_stderr = true
stdout_logfile = /usr/local/miserver/logs/materialscommons/mcbridgefsd.log
environment = HOME="/home/gtarcea",USER="gtarcea",MC_DOTENV_PATH="/home/gtarcea/workspace/src/github.com/materials-commons/materialscommons/.env"